# templates/nexus-wg-refresh.service.j2
[Unit]
Description=Refresh WireGuard AllowedIPs for Nexus (domain list)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecCondition=/usr/bin/systemctl is-active wg-quick@{{ nexus_wireguard_interface }}.service --quiet
ExecStart=/usr/local/sbin/nexus-wg-refresh-allowed-ips.sh
